# monero/services/monero-node-mainnet.service
[Unit]
Description=Monero Node Mainnet
Requires=docker.service
After=network-online.target
Wants=network-online.target

[Service]
User=root
WorkingDirectory=
ExecStart=/bin/bash compose up mainnet
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target

